# jintan-weekly-report-sop/deploy/jintan-weekly-report.service
[Unit]
Description=Jintan Phase II Weekly Report Generator
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/personal-agent-workspace/jintan-weekly-report-sop
EnvironmentFile=/opt/personal-agent-workspace/.env
ExecStart=/usr/bin/python3 /opt/personal-agent-workspace/jintan-weekly-report-sop/scripts/report_pipeline.py /opt/personal-agent-workspace/jintan-weekly-report-sop/config.yaml

# Auto-commit only the jintan weekly report project changes
ExecStartPost=/bin/bash -c 'cd /opt/personal-agent-workspace && git add --all /opt/personal-agent-workspace/jintan-weekly-report-sop/ && (git diff --cached --quiet || (git commit -m "auto: jintan weekly report update $(date +%%Y-%%m-%%d-%%H:%%M)" && git push origin $(git branch --show-current))) ; exit 0'

User=agent
Group=agent

# Capture output for debugging
StandardOutput=append:/var/log/jintan-weekly-report.log
StandardError=append:/var/log/jintan-weekly-report.log

[Install]
WantedBy=multi-user.target
